# ansible/roles/shibernate/templates/shibernate.timer.j2
[Unit]
Description=Run hibernate-wakeup cycle every 4 minutes

[Timer]
OnCalendar=*-*-* *:0/4:00
Persistent=true

[Install]
WantedBy=timers.target
