# ansible/roles/shibernate/templates/shibernate.service.j2
[Unit]
Description=Hibernate the system

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl hibernate
